# systemd/user/oh-hermes-auto-improve.service
[Unit]
Description=Generate oh-hermes self-improvement report
Documentation=https://github.com/NousResearch/hermes-agent

[Service]
Type=oneshot
ExecStart=%h/.local/bin/oh-hermes auto-improve
