# systemd/user/oh-hermes-secretary.timer
[Unit]
Description=Daily oh-hermes secretary briefing

[Timer]
OnCalendar=*-*-* 08:00
Persistent=true
RandomizedDelaySec=20m
Unit=oh-hermes-secretary.service

[Install]
WantedBy=timers.target
