# systemd/user/oh-hermes-secretary.service
[Unit]
Description=Generate oh-hermes daily secretary briefing
Documentation=https://github.com/NousResearch/hermes-agent

[Service]
Type=oneshot
Environment=OH_HERMES_HEALTH_TIMEOUT=20
ExecStart=%h/.local/bin/oh-hermes secretary brief
ExecStart=%h/.local/bin/oh-hermes secretary action plan
ExecStart=%h/.local/bin/oh-hermes secretary learn review
ExecStart=%h/.local/bin/oh-hermes secretary sweep
ExecStart=%h/.local/bin/oh-hermes secretary audit
